# crontabs
*/5 * * * * ~/current/workers/twitterSearch.py

