# supervisord.conf
[program:worker]
command = python /home/dotcloud/current/workers/tweetstreamworker.py